# 01-使用镜像link/service/dockerfile/supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/local/nginx/sbin/nginx

[program:php-fpm]
command=/usr/local/php/sbin/php-fpm
